# scripts/com.akb1.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  AKB1 Command Center — macOS LaunchAgent
  Per-user agent that runs scripts/autostart.sh on every login:
    1. Wait for Docker Desktop
    2. docker compose up -d
    3. Wait for /health
    4. open http://localhost:9000 in the default browser

  Install:
    ./scripts/install-autostart.sh

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.akb1.dashboard.plist
    rm ~/Library/LaunchAgents/com.akb1.dashboard.plist

  Logs:
    /tmp/akb1-launch.log  (stdout + stderr from the helper)
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.akb1.dashboard</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-lc</string>
      <string>"$HOME/Documents/Claude/Claude_Code/Projects/akb1-command-center/scripts/autostart.sh"</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- Retry if the first attempt exits non-zero (Docker still waking up etc.).
         Will throttle by ThrottleInterval seconds between retries. -->
    <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key>
      <false/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>300</integer>

    <key>StandardOutPath</key>
    <string>/tmp/akb1-launch.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/akb1-launch.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>LANG</key>
      <string>en_US.UTF-8</string>
    </dict>
  </dict>
</plist>
